# expo-fr-health.timer
[Unit]
Description=Expo-City Dubai FR Health Monitor Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
